# etc/systemd/system/hook-command.socket
[Unit]
Description=Hook Command socket activation

[Socket]
ListenStream=10292

[Install]
WantedBy=sockets.target
